# .config/aconfmgr/files/etc/cron.d/timeshift-boot
SHELL=/bin/bash
PATH=/usr/local/sbin:/usr/local/bin:/sbin:/bin:/usr/sbin:/usr/bin
MAILTO=""

@reboot root sleep 10m && timeshift --create --scripted --tags B
