# .config/systemd/user/swww@blur.service.d/override.conf
[Unit]
After=swww@wallpaper.service
Requires=swww@wallpaper.service
